# syncfail2band.service
[Unit]
Description=Syncfail2ban Service
After=network.target

[Service]
Type=simple
ExecStart=
Restart=on-failure
