# syncfail2band.service.template
# syncfail2band.service.template
[Unit]
Description=Syncfail2ban Service
After=network.target

[Service]
Type=simple
ExecStart=/etc/fail2ban/__init__.py
Restart=on-failure
